# supervisor_etc/stackalytics_dashboard.conf
[program:stackalytics_dashboard]
command=/usr/local/bin/stackalytics-dashboard --config-file /home/ken1ohmichi/commit-counter/stackalytics_etc/stackalytics.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log

